# resources/client/supervisord-chp.conf
[program:chp]
command=configurable-http-proxy --ip=0.0.0.0 --port=10000 --api-port=1999 --default-target='http://localhost:10002'
autostart=true
autorestart=true
priority=0